# com.quayside.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.quayside.pipeline</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/neilpeacock/Projects/quayside/scripts/run_if_needed.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/neilpeacock/Projects/quayside</string>
	<key>StartInterval</key>
	<integer>600</integer>
	<key>StandardOutPath</key>
	<string>/Users/neilpeacock/Projects/quayside/data/launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/neilpeacock/Projects/quayside/data/launchd.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin</string>
		<key>PYTHONPATH</key>
		<string>/Users/neilpeacock/Projects/quayside/src</string>
	</dict>
</dict>
</plist>
