# deploy/quayside-quality.timer
[Unit]
Description=Quayside data quality check — 3× daily weekdays

[Timer]
OnCalendar=Mon..Fri *-*-* 09:30:00
OnCalendar=Mon..Fri *-*-* 13:30:00
OnCalendar=Mon..Fri *-*-* 17:30:00
Persistent=true

[Install]
WantedBy=timers.target
