# deploy/quayside-quality.timer
[Unit]
Description=Quayside data quality check — backstop 3× daily weekdays (also runs after every scrape)

[Timer]
OnCalendar=Mon..Fri *-*-* 10:00:00
OnCalendar=Mon..Fri *-*-* 13:00:00
OnCalendar=Mon..Fri *-*-* 16:00:00
Persistent=true

[Install]
WantedBy=timers.target
